# supervisord.conf
[supervisord]
nodaemon=true

[program:filebeat]
command=filebeat -c /usr/src/app/filebeat.yml &

[program:controller]
command=ts-node controller.ts